# files/etc/hotplug.d/iface/80-nswa
#!/bin/sh

case "$ACTION" in
ifup)
	case "$INTERFACE" in
	wan)
		/lib/ranga/run-hooks "wanup"
		;;
	esac
	;;
esac
